# docker/dev/services.d/horizon/finish
#!/bin/sh
/usr/local/bin/php /var/www/html/artisan horizon:terminate
echo "Horizon stopped."
